# receiver/hot_plug.service.j2
[Unit]
Description=HDMI recovery

[Service]
Type=simple
RemainAfterExit=no
User={{ user }}
ExecStart=/bin/bash {{ user_home }}/recover_hdmi.sh

[Install]
WantedBy=multi-user.target
